# json_script.timer
[Unit]
Description=Run json_script.sh every second

[Timer]
OnBootSec=10s
OnUnitActiveSec=1s
AccuracySec=1s

[Install]
WantedBy=timers.target