# json_script.service
[Unit]
Description=Script to generate JSON every second

[Service]
ExecStart=/home/cti/APCUPS/json_script.sh